# systemd/seka-game.service
[Unit]
Description=Seka Card Game
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=seka
Group=seka
WorkingDirectory=/var/www/seka-game
Environment="PATH=/var/www/seka-game/venv/bin"
Environment="PYTHONPATH=/var/www/seka-game"
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/var/www/seka-game/venv/bin/python server.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=seka-game

[Install]
WantedBy=multi-user.target 